# rootfs/etc/s6-overlay/s6-rc.d/go-cqhttp/run
#!/usr/bin/with-contenv sh

set -ex && sleep 1 && /usr/local/bin/go-cqhttp -w /apps/go-cqhttp/conf.d || exit 1

exit 0
